# assets/config/supervisor/services-tpl.ini
[program:ss-server]
command=ss-server -c /etc/myconfig/ss/server.json -u
autostart=true
autorestart=true
startsecs=3

[program:v2ray-server]
command=/usr/local/bin/v2ray/v2ray --config=/etc/myconfig/v2ray/config.json
autostart=true
autorestart=true
startsecs=3

[program:brook-server]
command=/usr/local/bin/brook servers -l ":12089 --brookpw--"
autostart=true
autorestart=true
startsecs=3
